# components/example-notebook-servers/codeserver-python-rootless/s6/services.d/dockerd/finish
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "Docker daemon stopped"

# Clean up any remaining processes
pkill -f dockerd || true